# romm/rootfs/etc/services.d/romm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the RomM service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting RomM..."

# Declare variables
declare db_host
declare db_name
declare db_user
declare db_passwd
declare auth_secret_key
declare library_path
declare assets_path

# Get configuration values
db_host=$(bashio::config 'db_host')
db_name=$(bashio::config 'db_name')
db_user=$(bashio::config 'db_user')
db_passwd=$(bashio::config 'db_passwd')
auth_secret_key=$(bashio::config 'auth_secret_key')
library_path=$(bashio::config 'library_path')
assets_path=$(bashio::config 'assets_path')

# Check for MariaDB service if not configured manually
if [ -z "$db_host" ] && bashio::services.available 'mariadb'; then
    bashio::log.info "MariaDB service detected. Using MariaDB add-on for database."
    db_host=$(bashio::services 'mariadb' 'host')
    db_name=$(bashio::config 'db_name') # Keep from config if set, otherwise default
    db_user=$(bashio::services 'mariadb' 'username')
    db_passwd=$(bashio::services 'mariadb' 'password')
fi

# Export environment variables for RomM
if [ -n "$db_host" ]; then
    export DB_HOST="${db_host}"
    export DB_NAME="${db_name:-romm}"
    export DB_USER="${db_user}"
    export DB_PASSWD="${db_passwd}"
    export ROMM_DB_DRIVER="mariadb"
    bashio::log.info "Using database at ${db_host}"
else
    bashio::log.error "No database configured! RomM requires a MariaDB/MySQL or PostgreSQL database."
    bashio::log.error "Please configure the database settings or install the MariaDB add-on."
    bashio::exit.nok "Missing database configuration"
fi

# Set auth secret key
if [ -n "$auth_secret_key" ]; then
    export ROMM_AUTH_SECRET_KEY="${auth_secret_key}"
else
    bashio::log.warning "No auth secret key configured. Generating a random one..."
    export ROMM_AUTH_SECRET_KEY=$(python3 -c 'import secrets; print(secrets.token_hex(32))')
    bashio::log.info "Generated ROMM_AUTH_SECRET_KEY"
fi

# Export optional API keys
screenscraper_user=$(bashio::config 'screenscraper_user')
screenscraper_password=$(bashio::config 'screenscraper_password')
if [ -n "$screenscraper_user" ] && [ -n "$screenscraper_password" ]; then
    export SCREENSCRAPER_USER="${screenscraper_user}"
    export SCREENSCRAPER_PASSWORD="${screenscraper_password}"
    bashio::log.info "ScreenScraper configured"
elif [ -n "$screenscraper_user" ]; then
    bashio::log.warning "ScreenScraper user provided but password is missing"
fi

retroachievements_api_key=$(bashio::config 'retroachievements_api_key')
if [ -n "$retroachievements_api_key" ]; then
    export RETROACHIEVEMENTS_API_KEY="${retroachievements_api_key}"
    bashio::log.info "RetroAchievements configured"
fi

steamgriddb_api_key=$(bashio::config 'steamgriddb_api_key')
if [ -n "$steamgriddb_api_key" ]; then
    export STEAMGRIDDB_API_KEY="${steamgriddb_api_key}"
    bashio::log.info "SteamGridDB configured"
fi

# SSL Configuration
if bashio::config.true 'ssl'; then
    CERTFILE=$(bashio::config 'certfile')
    KEYFILE=$(bashio::config 'keyfile')
    
    if [ -f "/ssl/${CERTFILE}" ] && [ -f "/ssl/${KEYFILE}" ]; then
        bashio::log.info "Configuring SSL for RomM..."
        export ROMM_HTTPS_ENABLED="true"
        export ROMM_HTTPS_CERT_PATH="/ssl/${CERTFILE}"
        export ROMM_HTTPS_KEY_PATH="/ssl/${KEYFILE}"
    else
        bashio::log.error "SSL certificates not found at /ssl/${CERTFILE} or /ssl/${KEYFILE}"
        bashio::exit.nok "Missing SSL certificates"
    fi
fi

hasheous_api_enabled=$(bashio::config 'hasheous_api_enabled')
if [ -n "$hasheous_api_enabled" ]; then
    export HASHEOUS_API_ENABLED="${hasheous_api_enabled}"
else
    export HASHEOUS_API_ENABLED="true"
fi

# Set romm base path and user-configured paths
export ROMM_BASE_PATH="/romm"
export ROMM_LIBRARY_PATH="${library_path}"
export ROMM_ASSETS_PATH="${assets_path}"

bashio::log.info "Library path: ${ROMM_LIBRARY_PATH}"
bashio::log.info "Assets path: ${ROMM_ASSETS_PATH}"

# Create directories if they don't exist
mkdir -p "${library_path}"
mkdir -p "${assets_path}"
mkdir -p /romm/resources
mkdir -p /redis-data
mkdir -p /romm/config

# Start Redis server in the background
bashio::log.info "Starting Redis..."
redis-server --dir /redis-data --daemonize yes

# Give Redis time to start
sleep 2

bashio::log.info "Configuration complete. Starting RomM..."

# Execute romm's entrypoint script with the backend command
if [ -f "/docker-entrypoint.sh" ]; then
    exec /docker-entrypoint.sh python3 /backend/main.py
elif [ -f "/entrypoint.sh" ]; then
    exec /entrypoint.sh python3 /backend/main.py
elif [ -f "/app/entrypoint.sh" ]; then
    exec /app/entrypoint.sh python3 /backend/main.py
else
    # Direct fallback
    exec python3 /backend/main.py
fi
    bashio::log.info "Found RomM entrypoint at /entrypoint.sh"
    exec /entrypoint.sh
elif [ -f "/app/entrypoint.sh" ]; then
    bashio::log.info "Found RomM entrypoint at /app/entrypoint.sh"
    exec /app/entrypoint.sh
else
    bashio::log.error "Could not find RomM entrypoint script."
    # Try to find where it is and run it if found
    FOUND_ENTRY=$(find / -name "*entrypoint.sh" 2>/dev/null | head -n 1)
    if [ -n "$FOUND_ENTRY" ]; then
        bashio::log.info "Found entrypoint at $FOUND_ENTRY, executing..."
        exec "$FOUND_ENTRY"
    fi
    bashio::exit.nok "Entrypoint not found"
fi
