# romm/rootfs/etc/services.d/romm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the RomM service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting RomM..."

# Declare variables
declare db_host
declare db_name
declare db_user
declare db_passwd
declare auth_secret_key
declare library_path
declare assets_path

# Get configuration values
db_host=$(bashio::config 'db_host')
db_name=$(bashio::config 'db_name')
db_user=$(bashio::config 'db_user')
db_passwd=$(bashio::config 'db_passwd')
auth_secret_key=$(bashio::config 'auth_secret_key')
library_path=$(bashio::config 'library_path')
assets_path=$(bashio::config 'assets_path')

# Check for MariaDB service if not configured manually
if [ -z "$db_host" ] && bashio::services.available 'mariadb'; then
    bashio::log.info "MariaDB service detected. Using MariaDB add-on for database."
    db_host=$(bashio::services 'mariadb' 'host')
    db_name=$(bashio::config 'db_name') # Keep from config if set, otherwise default
    db_user=$(bashio::services 'mariadb' 'username')
    db_passwd=$(bashio::services 'mariadb' 'password')
fi

# Export environment variables for RomM
if [ -n "$db_host" ]; then
    export DB_HOST="${db_host}"
    export DB_NAME="${db_name:-romm}"
    export DB_USER="${db_user}"
    export DB_PASSWD="${db_passwd}"
    export ROMM_DB_DRIVER="mariadb"
    bashio::log.info "Using database at ${db_host}"
else
    bashio::log.error "No database configured! RomM requires a MariaDB/MySQL or PostgreSQL database."
    bashio::log.error "Please configure the database settings or install the MariaDB add-on."
    bashio::exit.nok "Missing database configuration"
fi

# Set auth secret key
if [ -n "$auth_secret_key" ]; then
    export ROMM_AUTH_SECRET_KEY="${auth_secret_key}"
else
    bashio::log.warning "No auth secret key configured. Generating a random one..."
    export ROMM_AUTH_SECRET_KEY=$(python3 -c 'import secrets; print(secrets.token_hex(32))')
    bashio::log.info "Generated ROMM_AUTH_SECRET_KEY"
fi

# Export optional API keys
screenscraper_user=$(bashio::config 'screenscraper_user')
screenscraper_password=$(bashio::config 'screenscraper_password')
if [ -n "$screenscraper_user" ] && [ -n "$screenscraper_password" ]; then
    export SCREENSCRAPER_USER="${screenscraper_user}"
    export SCREENSCRAPER_PASSWORD="${screenscraper_password}"
    bashio::log.info "ScreenScraper configured"
elif [ -n "$screenscraper_user" ]; then
    bashio::log.warning "ScreenScraper user provided but password is missing"
fi

retroachievements_api_key=$(bashio::config 'retroachievements_api_key')
if [ -n "$retroachievements_api_key" ]; then
    export RETROACHIEVEMENTS_API_KEY="${retroachievements_api_key}"
    bashio::log.info "RetroAchievements configured"
fi

steamgriddb_api_key=$(bashio::config 'steamgriddb_api_key')
if [ -n "$steamgriddb_api_key" ]; then
    export STEAMGRIDDB_API_KEY="${steamgriddb_api_key}"
    bashio::log.info "SteamGridDB configured"
fi

# SSL Configuration
if bashio::config.true 'ssl'; then
    CERTFILE=$(bashio::config 'certfile')
    KEYFILE=$(bashio::config 'keyfile')
    
    if [ -f "/ssl/${CERTFILE}" ] && [ -f "/ssl/${KEYFILE}" ]; then
        bashio::log.info "Configuring SSL for RomM..."
        export ROMM_HTTPS_ENABLED="true"
        export ROMM_HTTPS_CERT_PATH="/ssl/${CERTFILE}"
        export ROMM_HTTPS_KEY_PATH="/ssl/${KEYFILE}"
    else
        bashio::log.error "SSL certificates not found at /ssl/${CERTFILE} or /ssl/${KEYFILE}"
        bashio::exit.nok "Missing SSL certificates"
    fi
fi

hasheous_api_enabled=$(bashio::config 'hasheous_api_enabled')
if [ -n "$hasheous_api_enabled" ]; then
    export HASHEOUS_API_ENABLED="${hasheous_api_enabled}"
else
    export HASHEOUS_API_ENABLED="true"
fi

# Verify HASS mounts
if [ -d "/share" ]; then
    bashio::log.info "HASS /share mount detected."
else
    bashio::log.warning "HASS /share mount NOT detected. Local folder will be used."
fi

if [ -d "/media" ]; then
    bashio::log.info "HASS /media mount detected."
else
    bashio::log.warning "HASS /media mount NOT detected. Local folder will be used."
fi

# Set romm base path to a writable location
# This avoid permission issues with the built-in /romm directory
export ROMM_BASE_PATH="/data/romm"
mkdir -p "${ROMM_BASE_PATH}"

bashio::log.info "Library path: ${library_path}"
bashio::log.info "Assets path: ${assets_path}"

# Create the target directories on the HASS shares
mkdir -p "${library_path}"
mkdir -p "${assets_path}"

# Symlink them into the RomM root
ln -snf "${library_path}" "${ROMM_BASE_PATH}/library"
ln -snf "${assets_path}" "${ROMM_BASE_PATH}/assets"

# Resources are static and provided by the image - link them to our new root
ln -snf "/romm/resources" "${ROMM_BASE_PATH}/resources"

# Setup persistent config
mkdir -p "/data/config"
ln -snf "/data/config" "${ROMM_BASE_PATH}/config"

if [ ! -f "/data/config/config.yml" ]; then
    echo "{}" > "/data/config/config.yml"
fi

# Networking configuration
export ROMM_PORT="8080"
export HOST="0.0.0.0"

# Note: RomM v3.0+ uses an internal process manager (/init) that starts
# Valkey (Redis), Gunicorn (Backend), and Nginx (Frontend).
# We must call the official entrypoint and /init to ensure all services start.

# Internal directories needed by RomM/Nginx
# Use /data for everything writable to avoid permission issues in rootfs
mkdir -p /data/redis
mkdir -p /data/nginx/client_body
mkdir -p /data/nginx/proxy
mkdir -p /data/nginx/fastcgi
mkdir -p /data/nginx/uwsgi
mkdir -p /data/nginx/scgi
chmod -R 777 /data

# Create a custom Valkey configuration to use /data and disable failing persistence
mkdir -p /usr/local/etc/valkey
cat <<EOF > /usr/local/etc/valkey/valkey.conf
dir /data/redis
save ""
appendonly no
protected-mode no
bind 127.0.0.1
EOF

# Create an empty config.yml if it doesn't exist for RomM persistence
if [ ! -f "/data/config/config.yml" ]; then
    echo "{}" > "/data/config/config.yml"
fi
chmod 666 /data/config/config.yml

# We try to symlink the individual file if the directory is protected
ln -sf /data/config/config.yml /romm/config/config.yml 2>/dev/null || true

# Point RomM to use our /data directories for its own internal logic where possible
export REDIS_HOST="127.0.0.1" # This tells RomM to skip its internal watchdog and use ours
export REDIS_PORT="6379"

# Patch Nginx configuration to use /data for temporary files instead of /tmp
# and ensure it runs as root to avoid permission dropping issues.
if [ -f /etc/nginx/nginx.conf ]; then
    bashio::log.info "Patching Nginx configuration..."
    # Add user root to top of file
    sed -i '1iuser root;' /etc/nginx/nginx.conf
    # Patch paths
    sed -i 's|/tmp/client_body|/data/nginx/client_body|g' /etc/nginx/nginx.conf
    sed -i 's|/tmp/fastcgi|/data/nginx/fastcgi|g' /etc/nginx/nginx.conf
    sed -i 's|/tmp/proxy|/data/nginx/proxy|g' /etc/nginx/nginx.conf
    sed -i 's|/tmp/uwsgi|/data/nginx/uwsgi|g' /etc/nginx/nginx.conf
    sed -i 's|/tmp/scgi|/data/nginx/scgi|g' /etc/nginx/nginx.conf
fi

# Patch RomM's internal init script to avoid dropping permissions to 'romm' user
# as the Home Assistant environment doesn't allow the master process to chown directories.
if [ -f /usr/bin/romm-init ]; then
    bashio::log.info "Patching RomM init script to prevent permission dropping..."
    sed -i "s/nginx -g 'user romm;'/nginx/g" /usr/bin/romm-init
fi

# Deep-scanning size support: Patch RomM to report uncompressed size for archives.
# This fixes metadata matching for zipped ROMs on ScreenScraper and Hasheous.
ROMM_HANDLER="/romm/backend/handler/filesystem/roms_handler.py"
if [ -f "$ROMM_HANDLER" ]; then
    bashio::log.info "Patching RomM to support deep-scanning of zipped ROMs..."
    
    # Inject helper method into FSRomsHandler class
    # We use a unique marker (class definition) to avoid double-injecting on restart
    if ! grep -q "_get_deep_size" "$ROMM_HANDLER"; then
        sed -i "/class FSRomsHandler(FSHandler):/a \\
    def _get_deep_size(self, path):\\
        try:\\
            path_str = str(path).lower()\\
            if path_str.endswith('.zip'):\\
                with zipfile.ZipFile(path) as z: return max(f.file_size for f in z.infolist())\\
            if path_str.endswith(('.tar', '.tar.gz', '.tgz', '.tar.bz2', '.tar.xz')):\\
                import tarfile\\
                with tarfile.open(path) as t: return max(m.size for m in t.getmembers() if m.isfile())\\
            if path_str.endswith('.7z'):\\
                import subprocess\\
                res = subprocess.run(['7zz', 'l', '-slt', '-ba', str(path)], capture_output=True, text=True, check=True)\\
                sizes = [int(l.split('=')[1]) for l in res.stdout.splitlines() if l.startswith('Size = ')]\\
                if sizes: return max(sizes)\\
        except: pass\\
        return os.stat(path).st_size" "$ROMM_HANDLER"

        # Apply the helper call
        sed -i "s/file_size_bytes=os.stat(abs_file_path).st_size/file_size_bytes=self._get_deep_size(abs_file_path)/" "$ROMM_HANDLER"
    fi
fi

bashio::log.info "Starting local Valkey (Redis) on /data/redis..."
valkey-server /usr/local/etc/valkey/valkey.conf --daemonize yes

bashio::log.info "Configuration complete. Starting RomM services..."

# Execute the official entrypoint with the preserved official init script.
# Since REDIS_HOST is set, it will behave as if using an external Redis
# and skip the failing /redis-data logic.
chmod +x /docker-entrypoint.sh /usr/bin/romm-init
exec /docker-entrypoint.sh /usr/bin/romm-init
