# luanti/rootfs/etc/services.d/postgres/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

exec su-exec postgres postgres -D /data/postgres
